# lib/docker/chronos/chronos.conf
[program:mesos_master]
command='/usr/bin/mesos-init-wrapper master'
autostart=true
autorestart=true
stderr_logfile=/var/log/mesos-err.log
stdout_logfile=/var/log/mesos.log

[program:chronos]
command="/usr/local/chronos/bin/start-chronos.bash --master zk://localhost:2181/mesos --zk_hosts zk://localhost:2181/mesos --http_port 8081
autostart=true
autorestart=true
stderr_logfile=/var/log/mesos-err.log
stdout_logfile=/var/log/mesos.log